# crawler/image/etc/cron.d/stock
30 16 * * * root /testing_scripts/stock.sh > /testing_scripts/stock.txt 2>&1


