# web/image/etc/supervisor/supervisord.conf
[program:tornado-app]
priority=10
directory=/app
command=python app.py
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/tornado-app.log
redirect_stderr=true
